# bin/api_comments.service
[Unit]
Description=uWSGI instance to serve project api comments
After=network.target

[Service]
User=dev
Group=www-data
WorkingDirectory=/opt/api_comments
Environment="PATH=/opt/api_comments/venv/bin"
ExecStart=/opt/api_comments/venv/bin/uwsgi --ini uwsgi.ini

[Install]
WantedBy=multi-user.target